# 09/src/packtd.service
[Unit]
Description="Checking the logs"

[Service]
Type=simple
ExecStart=/usr/local/bin/packtd

[Install]
WantedBy=multi-user.target
